# debian/unipi-os-configurator.unipicheck.service
[Unit]
Description=Check Unipi Id fingerprint and call os-configurator
After=local-fs.target
Before=sysinit.target
DefaultDependencies=No

[Service]
Type=oneshot
ExecStartPre=/opt/unipi/tools/os-configurator --check --update
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
