# debian/unipi-os-configurator.unipicheck.service
[Unit]
Description=Check Unipi Id fingerprint and call os-configurator
After=local-fs.target
Before=basic.target
After=sysinit.target
Requires=sysinit.target
DefaultDependencies=No
ConditionVirtualization=No

[Service]
Type=oneshot
ExecStart=/usr/sbin/os-configurator --check --update
RemainAfterExit=yes

#[Install]
#WantedBy=unipi-os-configurator.unipicheck.path
